# lucee_rce/files/lucee.service
[Unit]
Description=Autostart Lucee
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/local/src/bin/
ExecStart=/usr/local/src/bin/startup.sh
Restart=on-abort
RestartSec=1

[Install]
WantedBy=multi-user.target