# apache_druid_rce/files/druid.service
[Unit]
Description=Apache Druid Micro Quickstart
After=network.target

[Service]
Type=forking
User=root
WorkingDirectory=/usr/local/apache-druid/
ExecStart=/usr/local/apache-druid/bin/start-micro-quickstart
Restart=on-abort
RestartSec=1

[Install]
WantedBy=multi-user.target